# src/systemd/lbw-monitor.service.in
[Unit]
Description=Load Balancer Monitor
After=network.target lbw-vrf.service
Requires=lbw-vrf.service

[Service]
LogNamespace=lbw
ExecStart=@pkglibexecdir@/lbw-monitor
Restart=on-failure
RestartPreventExitStatus=0 255

[Install]
WantedBy=multi-user.target
