# src/systemd/lbw-balancer.service.in
[Unit]
Description=Load Balancer for WAN
After=network.target lbw-vrf.service

[Service]
LogNamespace=lbw
ExecStart=@pkglibexecdir@/lbw-balancer
Restart=on-failure
RestartPreventExitStatus=0 255

[Install]
WantedBy=multi-user.target
