# src/systemd/lbw-balancer.service.in
[Unit]
Description=Load Balancer for WAN
After=network.target lbw-monitor.service
Requires=lbw-monitor.service

[Service]
ExecStart=@bindir@/lbw-balancer
Restart=on-failure
RestartPreventExitStatus=0 255

[Install]
WantedBy=multi-user.target
