# src/systemd/lbw-vrf.service.in
[Unit]
Description=Load Balancer VRF
After=network.target

[Service]
ExecStart=@bindir@/lbw-vrf
Restart=on-failure
RestartPreventExitStatus=0 255

[Install]
WantedBy=multi-user.target
